# pythonMQTT/solar.service
[Unit]
Description=solar

[Service]
ExecStart=/home/pi/solar.py

[Install]
WantedBy=multi-user.target
